# supervisord.conf
[supervisord]
nodaemon = false
[program:foo]
command= php bin/console app:mercure:relance
autostart=true
autorestart=true
startretries=3
startsecs = 0