# supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false
environment = KEY1="value1",KEY2="value2"
[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
prompt = supervisord
[program:foo]
command= php bin/console app:mercure:relance
autostart=true
autorestart=true
startretries=3
startsecs = 0